# debian/yeti-web.postinst
#!/bin/bash

user=yeti-web
services="yeti-web.service yeti-delayed-job.service yeti-prometheus.service yeti-libreoffice-headless.service"

case "$1" in 
	configure)
		#apply files permissions for new user
		chown -R $user:$user /opt/$user

		if [ -d /run/systemd/system ] && [ -x /usr/bin/deb-systemd-invoke ]; then
			deb-systemd-invoke restart rsyslog.service >/dev/null || true
		fi

		if [ -x /usr/bin/deb-systemd-helper ]; then
			deb-systemd-helper unmask ${services} >/dev/null || true
			for s in ${services}; do
				if deb-systemd-helper --quiet was-enabled ${s}; then
					deb-systemd-helper enable ${s} >/dev/null || true
				else
					deb-systemd-helper update-state ${s} >/dev/null || true
				fi
			done
		fi

        	rm -f /lib/systemd/system/yeti-cdr-billing.service || true

		if [ -d /run/systemd/system ]; then
			systemctl --system daemon-reload >/dev/null || true
			systemctl enable yeti-cdr-billing@cdr_billing
			systemctl disable yeti-libreoffice-headless
		fi

		if [ ! -f /opt/yeti-web/config/pgq_processors.yml ]; then
			cp -p /opt/yeti-web/config/pgq_processors.yml.distr /opt/yeti-web/config/pgq_processors.yml
		fi

                if [ ! -f /opt/yeti-web/config/yeti_web.yml ]; then
                        cp -p /opt/yeti-web/config/yeti_web.yml.distr /opt/yeti-web/config/yeti_web.yml
                fi
		
		#remove old init scripts
		rm -f /etc/init.d/yeti-cdr-billing /etc/init.d/yeti-web /etc/init.d/yeti-delayed-job || true

	;;
	abort-upgrade|abort-remove|abort-deconfigure)
		exit 0
	;;
	*)
		echo "postinst called with unknown argument \$1'" >&2
		exit 1
    ;;
esac

. /usr/share/debconf/confmodule

exit 0
